# power-button-gpio.service
[Unit]
Description=GPIO Power Button Monitoring Service
After=local-fs.target

[Service]
Type=simple
ExecStart=/bin/bash /path/to/power-button-gpio.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
