# power-button-gpio.service
[Unit]
Description=GPIO Power Button Monitoring Service
# Ensures the service starts after the local filesystem is mounted
After=local-fs.target

[Service]
Type=simple
ExecStart=/path/to/power-button-gpio.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
